# debian/ceilometer-collector.cron.daily
#!/bin/sh

if [ -x /usr/bin/ceilometer-expirer ] && [ -x /usr/bin/ceilometer-collector ] ; then
	su -s /bin/sh -c "exec /usr/bin/ceilometer-expirer" ceilometer
fi

exit 0
