# debian/ceilometer-alarm-notifier.upstart.in
description "ceilometer-alarm-notifier"
author "Thomas Goirand <zigo@debian.org>"

start on runlevel [2345]
stop on runlelvel [016]

chdir /var/run

pre-start script
	mkdir -p /var/run/ceilometer
	chown ceilometer:ceilometer /var/run/ceilometer

	mkdir -p /var/lock/ceilometer
	chown ceilometer:ceilometer /var/lock/ceilometer
end script

script
	[ -r /etc/default/openstack ] && . /etc/default/openstack
	[ -r /etc/default/$UPSTART_JOB ] && . /etc/default/$UPSTART_JOB

	[ "x$USE_SYSLOG" = "xyes" ] && DAEMON_ARGS="$DAEMON_ARGS --use-syslog"
	[ "x$USE_LOGFILE" != "xno" ] && DAEMON_ARGS="$DAEMON_ARGS --log-file /var/log/ceilometer/$UPSTART_JOB.log"

	exec su -s /bin/sh -c "exec ceilometer-alarm-notifier --config-file /etc/ceilometer/ceilometer.conf $DAEMON_ARGS" ceilometer
end script
